# docker/servers/config/supervisord.conf
; Supervisor config for the `servers` container: it runs both sconcur master
; servers (HTTP and socket) as two supervised programs in one container.
;
; Each program runs `sconcur-server start` in the foreground (the master itself
; supervises its own worker pool). autorestart=unexpected + exitcodes=0 means a
; clean stop (`sconcur-server stop` removes the state file → the master shuts down
; and exits 0) is left down, while a crash (non-zero / signal) is restarted.
;
; Runs as the unprivileged image user, so all runtime files live under /tmp.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:http-server]
command=php /sconcur/bin/sconcur-server start --configPath=/sconcur/config/sconcur.http-server.config.json
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=2
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:socket-server]
command=php /sconcur/bin/sconcur-server start --configPath=/sconcur/config/sconcur.socket-server.config.json
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=2
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
